# com.eba.fb-marketing-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eba.fb-marketing-agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>/Users/eba/code/auto-fb/fb_marketing_agent.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/eba/code/auto-fb</string>
    <key>StandardErrorPath</key>
    <string>/Users/eba/code/auto-fb/agent_error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/eba/code/auto-fb/agent_output.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
